# crontab
* * * * * cd /var/www/customers && /bin/bash backup.sh >> /dev/null 2>&1
